# debian/prerm
#!/bin/sh
# prerm script for vx-dga-c-google-drive
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <prerm> `remove'
#        * <old-prerm> `upgrade' <new-version>
#        * <new-prerm> `failed-upgrade' <old-version>
#        * <conflictor's-prerm> `remove' `in-favour' <package> <new-version>
#        * <deconfigured's-prerm> `deconfigure' `in-favour'
#          <package-being-installed> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package


case "$1" in
    remove)
	
    # Eliminamos el montaje Bind del/etc/fstab:
        FSTAB="/etc/fstab"
        MONTAJE="${DIRDRIVE} ${DIRBINDGRIVE} none bind 0 0"
        if ( grep ^${MONTAJE} "${FSTAB}" &> /dev/null ) ; then
           sed -i "/^${MONTAJE}/d" "${FSTAB}"
        fi
    # Matamos el proceso cliente actual:
        ## /sbin/start-stop-daemon --stop --name "vx-grive2"
        echo "Matamos el proceso de grive2 ..."
        pgrep -f vx-grive2 &> /dev/null && pkill -f vx-grive2

    # Lo desconfiguramos como servicio al inicio:
        if test -f "/etc/init.d/vx-grive2" ; then
			if /usr/sbin/update-rc.d vx-grive2 remove ; then
				echo "=> Se ha eliminado a Grive2 como un Servicio del inicio ..."
			fi
        fi
    ;;
    upgrade|deconfigure)
	;;
    failed-upgrade)
    ;;

    *)
        echo "prerm called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
